# systemd/telegram-proxy-daemon.service
[Unit]
Description=Telegram Auto-Proxy Daemon (zero-tap MTProto switcher)
Documentation=https://github.com/eaccesdev/MTPROTOBOT-Finn-
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=180s
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/home/null/Desktop/TelegramProxyBot

ExecStart=/bin/bash -c '\
    if [ -f /home/null/Desktop/TelegramProxyBot/venv/bin/python ]; then \
        exec /home/null/Desktop/TelegramProxyBot/venv/bin/python auto_proxy_daemon.py; \
    else \
        exec python3 auto_proxy_daemon.py; \
    fi'

Restart=on-failure
RestartSec=15s

# Allow time for first-time Telegram login prompt
TimeoutStartSec=120

StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-proxy-daemon

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
