# systemd/telegram-proxy-bot.service
[Unit]
Description=Telegram Proxy Manager Bot
Documentation=https://github.com/eaccesdev/MTPROTOBOT-Finn-
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=120s
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/home/null/Desktop/TelegramProxyBot

ExecStart=/bin/bash -c '\
    if [ -f /home/null/Desktop/TelegramProxyBot/venv/bin/python ]; then \
        exec /home/null/Desktop/TelegramProxyBot/venv/bin/python bot.py; \
    else \
        exec python3 bot.py; \
    fi'

Restart=on-failure
RestartSec=10s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-proxy-bot

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
